# packages/botlib/botlib-104.tar.gz/botlib-104/files/botd.service
[Unit]
Description=24/7 channel daemon
After=network-online.target
Wants=network-online.target

[Service]
User=botd
Group=botd
ExecStart=/usr/local/bin/botd

[Install]
WantedBy=multi-user.target
