# packages/genocide/genocide-6.tar.gz/genocide-6/files/genocide.service
[Unit]
Description=GENOCIDE - OTP-CR-117/19/001 - otp.informationdesk@icc-cpi.int - https://genocide.rtfd.io

After=network-online.target
Wants=network-online.target

[Service]
User=genocide
Group=genocide
ExecStart=/usr/local/bin/genocide mods=irc,rss -s

[Install]
WantedBy=multi-user.target
